# recto.service
[Unit]
Description=Gunicorn instance to serve Recto
After=network.target

[Service]
# Adjust User to your ubuntu user, e.g., 'ubuntu' or 'root' if you really must (not recommended)
User=www-data
Group=www-data
WorkingDirectory=/var/www/recto
Environment="PATH=/var/www/recto/venv/bin"
ExecStart=/var/www/recto/venv/bin/gunicorn --workers 3 --timeout 120 --bind unix:recto.sock -m 007 recto.wsgi:application
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
